# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:gunicorn]
command=/path/to/gunicorn main:app -c /path/to/gunicorn.conf.py
directory=/path/to/your/app
user=your_username
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/path/to/gunicorn.log

[program:watchdog]
command=python /path/to/your/watchdog_script.py
directory=/path/to/your/app
user=your_username
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/path/to/watchdog.log
